# worker/systemd/tradingagents-worker.service
[Unit]
Description=TradingAgents VPS worker daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tradingagents
Group=tradingagents
WorkingDirectory=/opt/tradingagents
EnvironmentFile=/etc/tradingagents-worker.env
ExecStart=/opt/tradingagents/.venv/bin/python -m worker.daemon
Restart=always
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/tradingagents /var/log/tradingagents

# Resource limits — LLM calls occasionally spike memory when parsing
# large payloads; 4 GB is comfortable for CX22.
MemoryMax=3G
CPUQuota=200%

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
